# scripts/mygest-cleanup.service
[Unit]
Description=MyGest Auto File Cleanup Service
After=network.target

[Service]
Type=simple
User=sandro
WorkingDirectory=/home/sandro/mygest/scripts
ExecStart=/usr/bin/python3 /home/sandro/mygest/scripts/auto_cleanup_watcher.py --watch-dir /home/sandro/upload-staging --delay 300 --interval 60
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/mygest-cleanup.log
StandardError=append:/var/log/mygest-cleanup-error.log

[Install]
WantedBy=multi-user.target
